# sddm-breeze-bing-wallpaper.timer
[Unit]
Description=Update sddm-breeze theme's wallpaper from bing 

[Timer]
OnUnitActiveSec=8h
Persistent=true

[Install]
WantedBy=timers.target