# sddm-breeze-bing-wallpaper.timer
[Unit]
Description=Update sddm-breeze theme's wallpaper from bing

[Timer]
OnCalendar=0/2:0
Persistent=true

[Install]
WantedBy=timers.target
